# data/systemd/centos/groonga-httpd.service.in
[Unit]
Description=Groonga HTTP Server
After=network.target
Conflicts=groonga-server-http.service

[Service]
Type=forking
PIDFile=@GROONGA_HTTPD_PID_PATH@
User=root
Group=root
EnvironmentFile=-/etc/sysconfig/groonga-httpd
ExecStartPre=-/usr/bin/mkdir -p /var/log/groonga/httpd
ExecStart=/usr/sbin/groonga-httpd
ExecStop=/usr/sbin/groonga-httpd -s stop

[Install]
WantedBy=multi-user.target
