# data/systemd/centos/groonga-httpd.service.in
[Unit]
Description=Groonga full text search HTTP server based on nginx
Documentation=http://groonga.org/docs/server/package.html#groonga-httpd
After=network.target
Conflicts=groonga-server-http.service

[Service]
Type=forking
PIDFile=@GROONGA_HTTPD_PID_PATH@
User=root
Group=root
EnvironmentFile=-/etc/sysconfig/groonga-httpd
ExecStartPre=-/usr/bin/mkdir -p /var/log/groonga/httpd
ExecStartPre=/usr/bin/rm -f @GROONGA_HTTPD_PID_PATH@
ExecStartPre=/usr/sbin/groonga-httpd -t
ExecStart=/usr/sbin/groonga-httpd
KillMode=process
KillSignal=SIGQUIT
PrivateTmp=true

[Install]
WantedBy=multi-user.target
