# packages/debian/groonga-server-gqtp.service
[Unit]
Description=Groonga full text search GQTP server
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
EnvironmentFile=-/etc/default/groonga-server-gqtp
User=groonga
Group=groonga
ExecStart=/usr/bin/groonga --protocol gqtp --bind-address $ADDRESS --log-path $LOG_PATH --query-log-path $QUERY_LOG_PATH --port $PORT -s $DATABASE

[Install]
WantedBy=multi-user.target
