# packages/debian/groonga-httpd.service
[Unit]
Description=Groonga full text search HTTP server
Wants=network-online.target
After=network-online.target
Conflicts=groonga.service

[Service]
Type=forking
User=root
Group=root
ExecStartPre=-/bin/mkdir -p /var/log/groonga/httpd
ExecStart=/usr/sbin/groonga-httpd
ExecStop=/usr/sbin/groonga-httpd -s stop

[Install]
WantedBy=multi-user.target
