# packages/debian/groonga-httpd.service
[Unit]
Description=Groonga full text search HTTP server
Documentation=http://groonga.org/docs/server/package.html#groonga-httpd
Wants=network-online.target
After=network-online.target

[Service]
Type=forking
User=root
Group=root
EnvironmentFile=-/etc/default/groonga-httpd
ExecStartPre=-/bin/mkdir -p /var/log/groonga/httpd
ExecStartPre=-/bin/mkdir -p /run/groonga
ExecStart=/usr/sbin/groonga-httpd
ExecStop=/usr/sbin/groonga-httpd -s stop

[Install]
WantedBy=multi-user.target
