# packages/debian/groonga-server-gqtp.service
[Unit]
Description=Groonga GQTP Server
Wants=network-online.target
After=network-online.target

[Service]
Type=forking
EnvironmentFile=-/etc/default/groonga-server-gqtp
User=groonga
Group=groonga
ExecStart=/usr/bin/groonga --protocol gqtp -d --pid-path /var/run/groonga/groonga-gqtp.pid --bind-address $ADDRESS --log-path $LOG_PATH --query-log-path $QUERY_LOG_PATH --port $PORT $DATABASE

[Install]
WantedBy=multi-user.target
